# hosting/ai.hermes.manager.service
[Unit]
Description=Hermes Manager
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash scripts/start-prod.sh
WorkingDirectory=/opt/hermes-manager
Environment=PORT=18470
Environment=NODE_ENV=production
StandardOutput=append:/opt/hermes-manager/runtime/logs/webapp.log
StandardError=append:/opt/hermes-manager/runtime/logs/webapp.error.log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
